# supervisord.conf
[supervisord]
nodemon=true

[program: pm2]
priority=1
command=pm2 start /app/process.yml

[program:sshd]
command=/etc/init.d/ssh start
autorestart=false
