# UI/supervisord.conf
[supervisord]
nodaemon=true

[program:apache_restart]
command=service apache2 restart
autorestart=false

[program:sshd]
command=/etc/init.d/ssh start
autorestart=false
